# packages/systemd/.config/systemd/user/kitty-save-session.service
[Unit]
Description=Save kitty session

[Service]
Type=oneshot
ExecStart=%h/.local/bin/kitty-save-session.sh
