# packages/systemd/.config/systemd/user/app-kitty@.service.d/exec.conf
[Service]
ExecStart=
ExecStart=%h/.local/bin/kitty-start.sh
